# packaging/systemd/nextdeploy-daemon.service
[Unit]
Description=NextDeploy Container Management Daemon
Documentation=https://github.com/aynaash/NextDeploy
After=network.target docker.service
Requires=docker.service
AssertPathExists=/etc/nextdeploy/config.yml

[Service]
Type=simple
# It's best practice to use an absolute path to the binary
ExecStart=/usr/local/bin/nextdeploy-daemon --config /etc/nextdeploy/config.yml --foreground

# Restart policy
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Security hardening - (CRITICAL: Create a dedicated user/group)
User=nextdeploy
Group=nextdeploy
# Consider using DynamicUser=yes as a modern, secure alternative
#DynamicUser=yes
#SupplementaryGroups=docker

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
# Remove this if the daemon needs to make outbound connections
#RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallArchitectures=native
LockPersonality=true
MemoryDenyWriteExecute=true
PrivateDevices=true
ProtectHostname=true
ProtectClock=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nextdeployd

[Install]
WantedBy=multi-user.target
